# layers/systemd/pymobiledevice3-tunneld.service
[Unit]
Description=pymobiledevice3 tunnel daemon for iOS developer services
After=network.target usbmuxd.service
Wants=usbmuxd.service

[Service]
Type=simple
Environment=PYTHONPATH=/home/jules/.local/lib/python3.14/site-packages
ExecStart=/home/jules/.local/bin/pymobiledevice3 remote tunneld
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
